# packaging/systemd/nexus.target
[Unit]
Description=Nexus Object Storage - All Services
After=network-online.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target
