# packaging/systemd/nexus-encrypt.service
[Unit]
Description=Nexus Encrypt Microservice
After=network-online.target
Wants=network-online.target
PartOf=nexus.target

[Service]
Type=simple
User=nexus
Group=nexus
ExecStart=/usr/local/bin/encrypt-service --address :50051
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=nexus.target
