# packaging/systemd/nexus-keystore.service
[Unit]
Description=Nexus KeyStore Microservice
After=network-online.target
Wants=network-online.target
PartOf=nexus.target

[Service]
Type=simple
User=nexus
Group=nexus
ExecStart=/usr/local/bin/keystore-service --address :50054
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/nexus

[Install]
WantedBy=nexus.target
